# qemu-constellation/agristar-bridge-hw.service
[Unit]
Description=Constellation Bridge (proto/COBS to Nova, REST/WS to UI)
After=network-online.target postgresql.service
Wants=network-online.target
# Production policy: don't enter a restart loop forever if startup is
# hopelessly broken. 10 restarts in 5 minutes ⇒ stop and surface the
# failure (`systemctl status`, alerting).
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=gellert
Group=dialout
WorkingDirectory=/home/gellert/Gellert/constellation/constellation-ui/server
ExecStart=/usr/bin/npx tsx src/index.ts
Restart=always
RestartSec=5
# Give the bridge a fair shutdown window — must beat the 8 s deadline
# inside index.ts shutdown() but stay under systemd defaults.
TimeoutStopSec=15
KillSignal=SIGTERM
# OS-level soft cap. Bridge holds pg pool (8) + WS clients + buffers;
# 512 MiB is ~10x steady-state. Crashing here is preferable to OOM-killing
# the whole rpi5.
MemoryMax=512M
# Don't let journald rate-limit production logs.
LogRateLimitIntervalSec=0
LogRateLimitBurst=0

# Hardening — Constellation runs untrusted UI code via WS so the
# blast-radius of a bridge compromise must be small. Each option below
# matches a real failure mode, not just paranoia.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
# Whitelist the directories the bridge legitimately writes to. Without
# these, ProtectHome=read-only blocks the log spool and ProtectSystem=full
# blocks the in-tree .sim-config cache (both observed as ENOENT/EROFS in
# 2026-04-29 hardening rollout).
ReadWritePaths=/home/gellert/.constellation /home/gellert/Gellert/constellation/constellation-ui/server/.sim-config /home/gellert/Gellert/upgrade /var/run/postgresql

# Real hardware: UART to Nova controller via GPIO14/15
Environment=SERIAL_PORT=/dev/ttyAMA0
Environment=PORT=9001
Environment=VFD_HOST=
Environment=NODE_ENV=production
Environment=PG_DSN=postgres:///paneldb?host=/var/run/postgresql
Environment=UPGRADE_DIR=/home/gellert/Gellert/upgrade
Environment=UPLOAD_DIR=/home/gellert/Gellert/upgrade/uploads

StandardOutput=journal
StandardError=journal
SyslogIdentifier=agristar-bridge

[Install]
WantedBy=multi-user.target
