# rpi5/agristar-pg-backup.timer
[Unit]
Description=Daily Constellation paneldb backup
Documentation=file:///home/gellert/Gellert/rpi5/agristar-pg-backup.sh

[Timer]
# 02:15 local — well clear of midnight log rotation and well before the
# bridge's 90 d retention sweep (runs on bridge boot + every 24 h from
# whenever it last booted).
OnCalendar=*-*-* 02:15:00
# If the rpi5 was off at 02:15, run on next boot.
Persistent=true
# Tiny jitter so a fleet of panels doesn't all hammer their disks at
# the exact same wall-clock instant.
RandomizedDelaySec=5min
Unit=agristar-pg-backup.service

[Install]
WantedBy=timers.target
